# systemd/system/etux-vpnserver-auth@etux-vpnserver.service.d/bind-to-default-service.conf
[Unit]
BindsTo=etux-vpnserver-auth.service
