# systemd/system/etux-vpnserver-sqluserportal@sqluserportal.service.d/bind-to-default-service.conf
[Unit]
BindsTo=etux-vpnserver-sqluserportal.service
